# wordpressMalScan.cron
#!/bin/sh
php /var/www/html/phpMalCodeScanner.php
